# examples/supervisord.conf
[program:spatial_index_api]
directory=/home/ubuntu/remote
command=gunicorn spatial_index_api:application
    --bind 0.0.0.0:{port} --workers 1
    --env ACCESS_KEY={access_key}
    --env USER_URI={user_uri}
autostart=true
autorestart=true
stopasgroup=true
stopsignal=QUIT

[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/home/ubuntu/supervisord.log            ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/home/ubuntu/supervisord.pid            ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
childlogdir=/home/ubuntu/supervisord/           ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
